# misc/crontab/storagebox
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
# 
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
# 
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
# 
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command

#Place the current date and time in the homelab dir so that when archives are taken its easily read
@daily date > /homelab/archive-date.txt

#Backup nextcloud's megasync directory to mega-xanthus
@monthly rclone sync local:/homelab/nextcloud/data/xanthus/files/megasync/ mega-xan:/ 
@weekly rclone cleanup mega-xan:/

#Use borgbackup to sync nextcloud locally and rclone to hetznerstoragebox
@daily borg create --list /mnt/borgback/nextcloud::archive-{hostname}-{now} /homelab/nextcloud/data/xanthus/files/
@weekly rclone sync local:/mnt/borgback/nextcloud/ hetznerstoragebox:/nextcloud-borgbackup-archive
@weekly rclone cleanup hetznerstoragebox:/

#Backup Homelab Dir to hetzner
@weekly rclone sync local:/homelab/ hetznerstoragebox:/homelab-mirror

#backup v's mc server
#@daily rclone sync gamebox:/homelab/viriditys-mc local:/home/temp && cd /home/temp && 7z a -t7z -m0=lzma -mx=9 -mfb=64 -md=32m -ms=on /homelab/backups/viriditys-mc/viriditys-mc_$(date +\%Y\%m\%d).7z * && rm -rf /home/temp

#Backup mc-xvr from oraclevps locally
@weekly rclone sync oraclevps:/home/ubuntu/mc-xvr local:/home/temp && cd /home/temp && 7z a -t7z -m0=lzma -mx=9 -mfb=64 -md=32m -ms=on /homelab/backups/mc-xvr/mc-xvr_$(date +\%Y\%m\%d).7z * && rm -rf /home/temp
