# crontab/storagebox
# Edit this file to introduce tasks to be run by cron.
#
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
#
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
#
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
#
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
#
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
#
# For more information see the manual pages of crontab(5) and cron(8)
#
# m h  dom mon dow   command

#Backup ampbox conan instances to local daily then push to uvps
@weekly rclone sync ampbox:/home/amp/.ampdata/instances/conansmp01/Backups local:/homelab/ampbox-backup/conan-exiled-lands #Pull from ampbox to local
@weekly rclone sync local:/homelab/ampbox-backup/conan-exiled-lands uvpsroot:/homelab/backup/conan/ #Push from local to uvps

#Backup ubuntubox mincraft servers to local daily
@weekly rclone sync uvpsroot:/home/amp/.ampdata/instances/mc-xvr01/Backups local:/homelab/ubuntubox-backup/mc-xvr #Pull from uvps to local
@weekly rclone sync uvpsroot:/home/amp/.ampdata/instances/vanilla-smp01/Backups local:/homelab/ubuntubox-backup/reddit-smp #Pull from uvps to local

#Back up my nextcloud to mega on server start and daily and ubuntubox
@weekly rclone sync local:/homelab/nextcloud/data/xanthus/files/megasync/ mega-xan:/ --log-level INFO --log-file=/home/xanthus/mega-rclone.log #Push from local to mega
@weekly rclone sync local:/homelab/nextcloud/data/xanthus/files/megasync/ uvpsroot:/home/xanthus --log-level INFO --log-file=/home/xanthus/uvps-rclone.log #Push fromn local to uvps
@monthly 7z a -mx=9 /homelab/nextcloud-archive/$(date +"\%Y-\%m-\%d").7z /homelab/nextcloud/data/xanthus/files/megasync/ > /homelab/nextcloud-archive/$(date +"\%Y-\%m-\%d").log
